# deploy/ha-addon/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/finish
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.warning "nginx stopped (exit code: ${1}), restarting..."
sleep 3
